# startup/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log  ; main log file
loglevel=info  ; log level

# Add TZ to the global environment
environment=TZ="%(ENV_TZ)s"

# crond removed - background tasks now handled by internal Rust scheduler

# Legacy startup/refresh tasks removed - now handled by internal Rust scheduler
# The Rust API now handles all background tasks internally with tokio-cron-scheduler

[program:client_api]
command=/usr/local/bin/pinepods-api
redirect_stderr=true
stdout_logfile=/var/log/supervisor/client_api.log
stderr_logfile=/var/log/supervisor/client_api.log
stdout_logfile_maxbytes=10000
stopwaitsecs=5

[program:gpodder_api]
command=bash -c 'export DB_USER="%(ENV_DB_USER)s"; export DB_HOST="%(ENV_DB_HOST)s"; export DB_PORT="%(ENV_DB_PORT)s"; export DB_NAME="%(ENV_DB_NAME)s"; export DB_PASSWORD="%(ENV_DB_PASSWORD)s"; export SERVER_PORT=8042; /usr/local/bin/gpodder-api'
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/gpodder_api.log
stderr_logfile=/var/log/supervisor/gpodder_api.log
stdout_logfile_maxbytes=10000
stopwaitsecs=10


[program:main_app]
command=nginx -g 'daemon off;'
redirect_stderr=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx_error.log
stdout_logfile_maxbytes=10000
stopwaitsecs=5
